# overlay/root/.plugin/bin/install-hacs
#!/usr/bin/env bash

## HACS is a custom integration for Home Assistant.
## See the HACS website for details: https://hacs.xyz/

. /etc/rc.subr && load_rc_config homeassistant

# shellcheck disable=SC2154
homeassistant_install_hacs() {
  local custom_components hacs zip
  custom_components="${1:-"${homeassistant_config_dir}"}/custom_components"
  hacs="${custom_components}/hacs"
  if [ -d "${hacs}" ]; then
    echo -e "\nHACS directory already exist, cleaning up..."
    rm -R "${hacs}"
  elif [ ! -d "${custom_components}" ]; then
    install -d -g "${homeassistant_group}" -m 775 -o "${homeassistant_user}" -- "${custom_components}"
  fi
  ## Download hacs.zip to a tmp_file and test, then extract to final location.
  zip="$(mktemp -t hacs)"
  wget -O "${zip}" https://github.com/hacs/integration/releases/latest/download/hacs.zip
  unzip -qt "${zip}" || { rm "${zip}" ; err 1 "failed to extract" ; }
  unzip -d "${hacs}" "${zip}" \
  && chown -R "${homeassistant_user}":"${homeassistant_group}" "${hacs}" \
  && echo -e "\nRemember to restart Home Assistant before you configure HACS\n https://hacs.xyz/docs/configuration/start \n"
  rm "${zip}"
}

homeassistant_install_hacs "${1}"
