# overlay/usr/local/etc/rc.d/appdaemon
#!/bin/sh

# PROVIDE: appdaemon
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# appdaemon_enable:	Set to YES to enable the appdaemon service.
#			Default: NO
# appdaemon_user:	The user account used to run the appdaemon daemon.
#			This is optional, however do not specifically set this to an
#			empty string as this will cause the daemon to run as root.
#			Default: USER_NAME
# appdaemon_group:	The group account used to run the appdaemon daemon.
#			This is optional, however do not specifically set this to an
#			empty string as this will cause the daemon to run with group wheel.
#			Default: USER_GROUP
# appdaemon_config_dir:	Directory where appdaemon config is located.
#			Default: ~/appdaemon/conf
#
# copy this file to /usr/local/etc/rc.d/appdaemon 
# chmod +x /usr/local/etc/rc.d/appdaemon
#
# sysrc appdaemon_enable=yes                # to enable
# sysrc appdaemon_config_dir="~/appdaemon"  # to change
# service appdaemon start                   # to start

. /etc/rc.subr
name=appdaemon
rcvar=${name}_enable

pidfile_child="/var/run/${name}.pid"
pidfile="/var/run/${name}_daemon.pid"
logfile="/var/log/${name}_daemon.log"

load_rc_config ${name}
: ${appdaemon_enable:="NO"}
: ${appdaemon_user:="hass"}
: ${appdaemon_group:="hass"}
: ${appdaemon_config_dir:="/home/hass/appdaemon/conf"}

command="/usr/sbin/daemon"

start_precmd=${name}_precmd
appdaemon_precmd()
{
        rc_flags="-f -o ${logfile} -P ${pidfile} -p ${pidfile_child} /srv/appdaemon/bin/appdaemon -c ${appdaemon_config_dir} ${rc_flags}"
        
    if [ ! -e "${pidfile_child}" ]; then
            install -g ${appdaemon_group} -o ${appdaemon_user} -- /dev/null "${pidfile_child}";
    fi

    if [ ! -e "${pidfile}" ]; then
            install -g ${appdaemon_group} -o ${appdaemon_user} -- /dev/null "${pidfile}";
    fi

    if [ ! -e "${logfile}" ]; then
        install -g ${appdaemon_group} -o ${appdaemon_user} -- /dev/null "${logfile}";
    fi
 
    if [ ! -d "${appdaemon_config_dir}" ]; then
        install -d -g ${appdaemon_group} -o ${appdaemon_user} -- "${appdaemon_config_dir}";
    fi
}

stop_postcmd=${name}_postcmd
appdaemon_postcmd()
{
    rm -f -- "${pidfile}"
    rm -f -- "${pidfile_child}"
}

run_rc_command "$1"
